# bongo.ini
[program:bongo]
command=/usr/bin/node /var/www/virtual/tmwnd/bongo.tmwnd.de/app.js